# rtpp/install/.svn/text-base/rtppmon.svn-base
#!/bin/bash
#
# Startup script for the  monitor
#
# chkconfig: - 85 15
# description:A symmetric RTP proxy 
#
# processname: rtpp 
# pidfile: /var/run/rtpp.pid
#
#
# created by boll.tan: ShenZhen keepc Digital Ltd.
#

# Source function library.
. /etc/rc.d/init.d/functions

CONFIG=/root/.rtpp/profile/rtpp.conf
LOGFILES=/var/log/rtpp/rtpp.log
LOCKFILE=/var/lock/subsys/rtpp
PID=/var/run/rtpp.pid
MODULE=/opt/nf/nf_fwd.ko

# Source networking configuration and check that networking is up.
if [ -f /etc/sysconfig/network ] ; then
        . /etc/sysconfig/network
        [ ${NETWORKING} = "no" ] && exit 0
fi


[ -x /usr/local/sbin/rtpp ] || exit 0


prog="rtpp Monitor"
RETVAL=0

start() {
	echo $"checking nf_fwd module... "
	if [ -f ${MODULE} ]; then
		mod=`lsmod | grep nf_fwd`
		if [ -n "$mod" ]; then
			echo  "nf_fwd already runing."
		else
			echo "start nf_fwd module."
			/sbin/insmod $MODULE
		fi
	fi

        echo -n $"Starting $prog: " 
        [ -f $LOGFILES ] && mv -f $LOGFILES $LOGFILES.bak
        cd /usr/local/sbin/
        dtach -n /tmp/dtach-rtpp /usr/local/sbin/rtpp  
        RETVAL=$?
        /etc/rc.d/init.d/rtppmon status
        [ $RETVAL = 0 ] && touch $LOCKFILE && success
        echo
        return $RETVAL
}

stop() {
        echo -n $"Stopping $prog: "
        killproc rtpp 
        RETVAL=$?
        [ $RETVAL = 0 ] && rm -f $LOCKFILE $PID
	rm -f /tmp/dtach-rtpp
        echo
	sleep 1
        return $RETVAL
}

stopall() {
	echo -n $"Stoping $prog and remove nf_fwd from kernel"	
	killproc rtpp 
	RETVAL=$?
	[ $RETVAL = 0 ] && rm -f $LOCKFILE $PID
	rm -f /tmp/dtach-rtpp
	echo
	sleep 1
	if [ -f ${MODULE} ]; then
                mod=`lsmod | grep nf_fwd`
                if [ -n "$mod" ]; then
			/sbin/rmmod nf_fwd
                        echo  "remove nf_fwd"
                else
                        echo "nf_fwd is not running"
                fi
        fi
}

# See how we were called.
case "$1" in
    start)
        [ -f $LOCKFILE ] || start
        ;;
    stop)
        stop
        ;;
    stopall)
	stopall
	;;
    restart)
        stop
        start
        ;;
    reload)
        kill -HUP `cat $PID`
        RETVAL=$?
        ;;
    status)
        status rtpp 
        RETVAL=$?
        ;;
    *)
        echo $"Usage: rtppmon {start|stop|stopall|restart|reload|status}"
        exit 1
esac

exit $RETVAL

